# systemd/imx-web-dashboard.service
[Unit]
Description=IMX500 Web Dashboard
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/pi/tool_for_imx
ExecStart=/home/pi/venv/bin/python /home/pi/tool_for_imx/web_dashboard/server.py --host 0.0.0.0 --port 8080 --bbox-host 127.0.0.1 --bbox-port 5005 --no-video --mjpeg-host 127.0.0.1 --mjpeg-port 8081
Restart=always
RestartSec=2
User=pi
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=-/etc/default/imx-web-dashboard

[Install]
WantedBy=multi-user.target
